# systemd/user/mail_live.timer
[Unit]
Description=Live Mail Synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=1h
Unit=mail_live.service

[Install]
WantedBy=timers.target
