# systemd/user/mail_gmail.timer
[Unit]
Description=Gmail Daily Synchronisation

[Timer]
OnStartupSec=2m
OnUnitActiveSec=8h
Unit=mail_gmail.service

[Install]
WantedBy=timers.target
