# systemd/user/mail_codex.timer
[Unit]
Description=Mail Synchronization Timer

[Timer]
OnStartupSec=10s
OnUnitActiveSec=1m
Unit=mail_codex.service

[Install]
WantedBy=timers.target
