# systemd/user/mail_uni.timer
[Unit]
Description=Uni Mail Synchronization Timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=1m
Unit=mail_uni.service

[Install]
WantedBy=timers.target
