# systemd/user/update.service
[Unit]
Description=Update prompt

[Service]
ExecStart=/home/gregor/dotfiles/scripts/updateClock.sh

[Install]
WantedBy=default.target
