# systemd/user/update.timer
[Unit]
Description=Prompt upgrades

[Timer]
OnCalendar=hourly
Persistent=true
Unit=update.service

[Install]
WantedBy=timers.target
